# admin-web/deploy/copimine-minecraft-discord-bridge.service
[Unit]
Description=CopiMine Minecraft to Discord Bridge
After=network.target copimine-admin.service copimine-minecraft.service
Wants=copimine-admin.service copimine-minecraft.service

[Service]
Type=simple
User=copimine
WorkingDirectory=/opt/copimine/admin-web
EnvironmentFile=/opt/copimine/admin-web/.env
Environment=COPIMINE_ENV_FILE=/opt/copimine/admin-web/.env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/copimine/admin-web/backend:/opt/copimine/admin-web
ExecStart=/opt/copimine/admin-web/.venv/bin/python -m backend.minecraft_discord_bridge
Restart=always
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
