# admin-web/deploy/copimine-discord-bot.service
[Unit]
Description=CopiMine Discord Bot
After=network.target copimine-admin.service
Wants=copimine-admin.service

[Service]
Type=simple
User=qwerty
WorkingDirectory=/opt/copimine/admin-web
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/copimine/admin-web/.venv/bin/python -m backend.discord_bot
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
