# admin-web/deploy/copimine-discord-bot.service
[Unit]
Description=CopiMine Discord Bot
After=network.target copimine-admin.service
Wants=copimine-admin.service

[Service]
Type=simple
User=copimine
WorkingDirectory=/opt/copimine/admin-web
EnvironmentFile=/opt/copimine/admin-web/.env
Environment=COPIMINE_ENV_FILE=/opt/copimine/admin-web/.env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/copimine/admin-web/backend:/opt/copimine/admin-web
ExecStart=/opt/copimine/admin-web/.venv/bin/python -m backend.discord_bot
Restart=always
RestartSec=300
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
